# packages/service/home-assistant/system.d/hassio-apparmor.service
[Unit]
Description=Hass.io AppArmor
Wants=hassio-supervisor.service
Before=docker-ce.service hassio-supervisor.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/sbin/hassio-apparmor

[Install]
WantedBy=multi-user.target
